# systemd/sshwarma.service
[Unit]
Description=sshwarma SSH MUD for humans and models
Documentation=https://github.com/tobert/sshwarma

[Service]
Type=simple
ExecStart=%h/src/sshwarma/target/release/sshwarma
Restart=always
RestartSec=5
TimeoutStopSec=30

# Paths are now XDG-based:
#   Data:   ~/.local/share/sshwarma/ (db, host_key)
#   Config: ~/.config/sshwarma/ (models.toml, lua scripts)
# Override with: SSHWARMA_DB, SSHWARMA_HOST_KEY, SSHWARMA_MODELS_CONFIG
Environment=RUST_LOG=sshwarma=info

# OpenTelemetry OTLP export (uncomment to enable)
# Traces and metrics are exported to the configured endpoint
#Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://localhost:4317
#Environment=OTEL_SERVICE_NAME=sshwarma

[Install]
WantedBy=default.target
